# daemon/cups-locald-systemd.service.in
[Unit]
Description=CUPS Local Spooler

[Service]
Type=dbus
BusName=org.openprinting.cups-locald
ExecStart=@sbindir@/cups-locald
